# src/systemd/pdns-dilithium.service
[Unit]
Description=PowerDNS Authoritative Server (Dilithium Port 5305)
After=network.target mysql.service

[Service]
Type=simple
ExecStart=/usr/sbin/pdns_server --config-dir=/etc/powerdns --config-name=pdns-dilithium
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
