# src/configs/systemd/pdns.service
[Unit]
Description=PowerDNS Authoritative Server (Port 5300)
After=network.target mysql.service

[Service]
Type=simple
ExecStart=/usr/sbin/pdns_server --config-dir=/etc/powerdns --config-name=pdns --daemon=no --guardian=no
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
